# next-ai-draw-io.service
[Unit]
Description=Next.js AI Draw.io Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/next-ai-draw-io/.next/standalone
Environment=NODE_ENV=production
Environment=PORT=6002
Environment=HOSTNAME=0.0.0.0
Environment=USE_HTTP=true
Environment=PUBLIC_URL=http://43.128.81.143:6002
Environment=NEXT_TELEMETRY_DISABLED=1

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/opt/next-ai-draw-io

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=next-ai-draw-io

# Start command (will be updated by deploy.sh with correct paths)
# Default to /usr/bin/nodejs (Debian/Ubuntu standard) or /usr/bin/node
ExecStart=/usr/bin/nodejs /opt/next-ai-draw-io/.next/standalone/server.js

[Install]
WantedBy=multi-user.target

